# dot_config/systemd/user/toolbox-autoupdate.timer
[Unit]
Description=Periodic auto-update of Fedora Toolbx container

[Timer]
OnStartupSec=1m
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
